# roles/counter-panda/files/counter-panda.conf
description "counter-panda project"
author  "gmoshiko"

start on runlevel [234]
stop on runlevel [0156]

chdir /tmp/counter-panda
exec /tmp/counter-panda/counter-panda.py 
respawn
