# roles/gify-panda/files/gify-panda.conf
description "gify panda project"
author  "gmoshiko"

start on runlevel [234]
stop on runlevel [0156]

chdir /tmp/gify-panda
exec /tmp/gify-panda/gify-panda.py
respawn
